# autostart/kanshi.desktop
[Desktop Entry]
Type=Application
Exec=/home/martijnvm/bin/kanshi
Hidden=false
NoDisplay=false
Name=Kanshi Autostart
