# autostart/swayrd.desktop
[Desktop Entry]
Type=Application
Exec=/home/martijnvm/bin/swayrd
Hidden=false
NoDisplay=false
Name=Swayrd Autostart
